# conf/zkrun3
#!/bin/sh

# chkconfig: - 75 50
# description: Starts and stops Zookeeper

cd /opt/zookeeper3
bin/zkServer.sh $1